# apps/pi-controller/systemd/cubism-controller.service
[Unit]
Description=Cubism input controller (volume-knob module selector)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust the path if you cloned the repo elsewhere.
WorkingDirectory=/home/pi/Documents/Projects/cubism/apps/pi-controller

# Node lives wherever your install method put it. Pick ONE of the two
# ExecStart lines below.
#
# Option A (recommended for nvm users): source nvm explicitly so the unit
# keeps working after `nvm install`/`nvm use` swaps the active version.
ExecStart=/bin/bash -c 'export NVM_DIR="/home/pi/.nvm"; . "$NVM_DIR/nvm.sh"; exec node /home/pi/Documents/Projects/cubism/apps/pi-controller/dist/index.js'
#
# Option B: hardcode the absolute path. Faster startup, but you must update
# this line whenever you `nvm install` a new Node version.
# Run `which node` as the `pi` user to confirm the path.
# ExecStart=/home/pi/.nvm/versions/node/v24.16.0/bin/node /home/pi/Documents/Projects/cubism/apps/pi-controller/dist/index.js

# Replace <desktop-host> with your desktop's LAN address (where `pnpm dev`
# or `pnpm start` runs the Cubism socket server on port 3000).
Environment=CUBISM_SERVER_URL=http://<desktop-host>:3000
Environment=CUBISM_DEVICE_ID=pi-holo-001
Environment=CUBISM_USER_ID=demo-user

Restart=always
RestartSec=2
TimeoutStopSec=5
KillMode=mixed
User=pi
SupplementaryGroups=input

# Uncomment to dump every key event to the journal for debugging.
# Environment=CUBISM_CONTROLLER_DEBUG=1

# Uncomment to point at specific event devices (comma-separated) and skip
# VID:PID auto-detection.
# Environment=CUBISM_CONTROLLER_DEVICES=/dev/input/event9,/dev/input/event10,/dev/input/event11

# Uncomment to remap which key codes count as next/prev. Defaults are
# KEY_VOLUMEUP=115 and KEY_VOLUMEDOWN=114.
# Environment=CUBISM_CONTROLLER_NEXT_KEY=115
# Environment=CUBISM_CONTROLLER_PREV_KEY=114

[Install]
WantedBy=multi-user.target
